# script_install/celerybeat_master.conf
; ================================
;  celery beat supervisor example
; ================================

[program:celerybeat]
; Set full path to celery program if using virtualenv
command=celery beat -A datarun -s /home/datarun/celery_info/celerybeat-schedule --loglevel=INFO --pidfile="/home/datarun/celery_info/celerybeat.pid"

; remove the -A myapp argument if you are not using an app instance

directory=/home/datarun
user=celery
numprocs=1
stdout_logfile=/home/datarun/celery_info/beat.log
stderr_logfile=/home/datarun/celery_info/beat.log
autostart=true
autorestart=true
startsecs=10

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
